# build/lua/additional/common/qos/etc/hotplug.d/iface/10-bqos
#!/bin/sh
# Based on https://github.com/skyformat99/eqos/blob/master/files/eqos.hotplug

[ "$ACTION" = "ifup" ] || exit 0
[ "$INTERFACE" = "lan" ] || exit 0

/etc/init.d/bqos start